# online-siteaccess.service
[Unit]
Description=Online SiteAccess Chat Server
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/online.siteaccess.ru/apps/server
EnvironmentFile=/var/www/online.siteaccess.ru/apps/server/.env
ExecStart=/usr/bin/node /var/www/online.siteaccess.ru/apps/server/dist/main.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=online-siteaccess

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
